# recipes-apps/print-ip/print-ip/print-ip.timer
[Unit]
Description=Timer service to update the IP on OLED each 10s

[Timer]
OnUnitActiveSec=10s
OnBootSec=10s

[Install]
WantedBy=timers.target
